# rootfs/etc/s6-overlay/s6-rc.d/svc-openbox/run
#!/command/with-contenv bash

set -euo pipefail

# Start Openbox
exec s6-setuidgid "$(id -nu "${PUID}")" openbox
